# scripts/quickstats.crontab
MAILTO=""
*/5 * * * *	/yoda/bin/quickstats
